# ans/upgrade/cfgdb/v2.6.1/files/cfgdb-disco.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='cfgdb-disco'>
    <service name='application/cfgdb-disco' type='service' version='0'>
        <create_default_instance enabled='true' />
        <single_instance />
        <dependency name='network' grouping='require_all' restart_on='error' type='service'>
            <service_fmri value='svc:/milestone/network:default' />
        </dependency>
        <dependency name='filesystem' grouping='require_all' restart_on='error' type='service'>
            <service_fmri value='svc:/system/filesystem/local' />
        </dependency>
        <method_context>
            <method_credential group='root' user='root' />  <!-- required for query_cfgdb -->
            <method_environment>
                <envvar name='PATH' value='/usr/local/sbin:/usr/local/bin:/opt/local/sbin:/opt/local/bin:/usr/sbin:/usr/bin:/sbin' />
            </method_environment>
        </method_context>
        <exec_method name='start' type='method' exec='/usr/local/bin/cfgdb_disco_server.py &amp;' timeout_seconds='60' />
        <exec_method name='stop' type='method' exec=':kill' timeout_seconds='10' />
        <property_group name='startd' type='framework'>
            <propval name='duration' type='astring' value='contract' />
            <propval name='ignore_error' type='astring' value='core,signal' />
        </property_group>
        <template>
            <common_name>
                <loctext xml:lang='C'>Danube Cloud cfgdb discovery service</loctext>
            </common_name>
        </template>
    </service>
</service_bundle>